# etc/services.d/user/run
#!/usr/bin/with-contenv sh

echo "Starting Python App"
exec s6-setuidgid user /bin/sh -c 'python -u /app/server.py'
